# deployment/cms-hub-worker.conf
[program:market-worker]
process_name=%(program_name)s
directory={{APP_ROOT}}
command={{PHP_BIN}} {{APP_ROOT}}/worker.php
environment=WORKER_HEARTBEAT_FILE="{{APP_ROOT}}/storage/logs/worker-heartbeat",WORKER_SLEEP_SECONDS="3",WORKER_MAX_JOBS="500",WORKER_MAX_RUNTIME="3600"
autostart=true
autorestart=true
startsecs=5
user={{APP_USER}}
numprocs=1
redirect_stderr=true
stdout_logfile={{APP_ROOT}}/storage/logs/worker.log
stopwaitsecs=3600
